# cronjobs/crontab
0 0 * * * /root/bin/housekeeping.sh 2>&1 1>/dev/null